# api/etc/systemd/system/db-backup.service
[Unit]
Description=PostgreSQL Backup
Requires=postgresql.service

[Service]
User=postgres
ExecStart=/home/api/backup.sh
Type=oneshot
